# assets/supervisord.ini
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:xorg]
command=/usr/bin/Xorg
autostart = true
priority = 0

[program:webchela]
command=/home/user/webchela/venv/bin/webchela
autostart = true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
